# infra/installer/systemd/laia-ui-server.service.tmpl
[Unit]
Description=LAIA ARCH console (operator-only web UI on loopback/Tailscale)
# No dependency on laia-gateway: the console talks to its own FastAPI backend;
# the old Wants= chained it to a unit that no longer runs (v0.11 model).
After=network.target
X-LaiaPathDeps=ui_server,ui_backend,laia_venv

[Service]
Type=simple
User=${LAIA_USER}
Group=${LAIA_USER}
# Path registry — defines LAIA_UI_SERVER, LAIA_UI_BACKEND, LAIA_VENV, etc.
# Layout v2: laia-pathd writes .env.paths to the ARCH runtime home, not LAIA_HOME.
EnvironmentFile=-/srv/laia/arch/.env.paths
WorkingDirectory=${LAIA_INSTALL_PREFIX}/.laia-core/laia-ui-server/backend
Environment=HOME=${LAIA_USER_HOME}
Environment=LAIA_HOME=${LAIA_HOME}
Environment=LAIA_UI_HOST=127.0.0.1
# Operator console: LOOPBACK ONLY (users must never reach ARCH —
# LAIA_ECOSYSTEM). Remote operator access goes through Tailscale on the
# host, never a laia-edge vhost. The old --host 0.0.0.0 exposed
# unauthenticated PTY terminals to the LAN.
ExecStart=${LAIA_INSTALL_PREFIX}/.laia-core/venv/bin/python -m uvicorn main:app --host 127.0.0.1 --port 8077
Restart=on-failure
RestartSec=3
# Deliberately NO ProtectSystem/NoNewPrivileges sandbox: the console spawns
# operator PTY terminals (Command Center) and bwrap — those lockdowns break
# them (and the snap pitfall 2026-06-03 if a terminal runs lxc).
# Containment = User non-root + loopback bind.

[Install]
WantedBy=multi-user.target
