# services/laia-agent-runtime/systemd/laia-agent.service
[Unit]
Description=LAIA personal agent runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=laia-agent
Group=laia-agent
WorkingDirectory=/opt/laia/agent
Environment=PYTHONPATH=/opt/laia/agent/src
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/laia/runtime/venv/bin/python -m laia_agent
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
